# conf/supervisord-tarantool.conf
[supervisord]
nodaemon=true
pidfile=/run/supervisor/pid

[unix_http_server]
file=//run/supervisor/sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/sock

[program:cron]
command=/usr/sbin/cron -f
autorestart=true
startretries=4294967295

[program:heartbeat]
command=/usr/share/wallarm-common/synccloud
environment=WALLARM_SYNCNODE=no
autorestart=true
startretries=4294967295

[program:tarantool]
command=/usr/bin/tarantool /usr/share/wallarm-tarantool/init.lua
environment=HOST=0.0.0.0
user=wallarm
group=wallarm
autorestart=true
startretries=4294967295
