# conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/wallarm/supervisord.log
pidfile=/run/supervisor/pid

[unix_http_server]
file=//run/supervisor/sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/sock

[program:collectd]
command=/usr/sbin/collectd -f -C /etc/collectd/wallarm-collectd.conf
autorestart=true
startretries=4294967295
redirect_stderr=true
stdout_logfile=/var/log/wallarm/collectd.log

[program:cron]
command=/usr/sbin/cron -f
autorestart=true
startretries=4294967295
redirect_stderr=true
stdout_logfile=/var/log/wallarm/cron.log

[program:tarantool]
command=/usr/bin/wtarantool /usr/share/wallarm-tarantool/init.lua
user=wallarm
group=wallarm
autorestart=true
startretries=4294967295
redirect_stderr=true
stdout_logfile=/var/log/wallarm/tarantool-out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
startretries=4294967295
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stderr_logfile=/dev/stdout

[program:appstructure]
command=/usr/bin/wallarm-appstructure-client --config-file /etc/wallarm/node.yaml
user=wallarm
group=wallarm
autorestart=true
startretries=4294967295
redirect_stderr=true
stdout_logfile=/var/log/wallarm/appstructure.log

[program:addnode_loop]
command=/usr/local/bin/addnode_loop
autorestart=false
startsecs=0
redirect_stderr=true
stdout_logfile=/var/log/wallarm/addnode_loop.log